# trees/1.5/etc/systemd/system/xsce-onboot.service
[Unit]
Description=Small script to clean up offline install, set up networks,etc
Requires=network.target
After=network.target

[Service]
Type=idle
ExecStart=/usr/bin/xs-onboot

[Install]
WantedBy=multi-user.target
